# config/systemd/user/battery.service
[Unit]
Description=Check battery

[Service]
Type=simple
ExecStart=/home/tsokkur/.local/bin/cron/battery.sh

