# config/systemd/user/duplicity-backup.service
[Unit]
Description=Run duplicity backup

[Service]
Type=oneshot
ExecStart=/bin/bash -c '~/.local/bin/bak -b && notify-send backup'
